# score-mcp-server/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/app.log
logfile_maxbytes=50MB
logfile_backups=10
pidfile=/var/run/supervisord.pid
# Redirect all supervisord output to the same log file
stdout_logfile=/var/log/supervisor/app.log
stderr_logfile=/var/log/supervisor/app.log

[program:uvicorn]
command=/app/start.sh
directory=/app
autostart=true
autorestart=true
# Redirect stderr to stdout so all output (stdout and stderr) goes to stdout
redirect_stderr=true
# Send uvicorn logs to stdout for Docker log aggregation
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
# PYTHONUNBUFFERED="1" disables Python's output buffering, ensuring logs appear immediately
# Without this, Python buffers stdout/stderr when not connected to a TTY (common in Docker),
# which can delay log output until the buffer fills or the process exits. This is critical
# for real-time log monitoring, debugging, and Docker log aggregation tools.
environment=PYTHONUNBUFFERED="1"
# Pass through all environment variables from the container
# The start.sh script will load .env file and export all variables

